# scorpio.conf
[program:scorpio]
command = gunicorn scorpio.wsgi_test -b 127.0.0.1:8023
directory = /root/scorpio/scorpio
autostart = true
autorestart = true

